# systemd/smart-mirror-ir.timer
[Unit]
Description=Run Smart Mirror IR mirror update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
